# Z/supervisord.conf
[program:node]
command = node worker.js
directory = /home/dotcloud/current

